# config/init/systemd/lxcfs.service.in
[Unit]
Description=FUSE filesystem for LXC
ConditionVirtualization=!container
Before=lxc.service
Documentation=man:lxcfs(1)

[Service]
ExecStart=/usr/bin/lxcfs {{LXCFSTARGETDIR}}
KillMode=process
Restart=on-failure
ExecStopPost=-/bin/fusermount -u {{LXCFSTARGETDIR}}
Delegate=yes
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
